# config/webmin
#!/sbin/openrc-run
#WEBMIN=/etc/webmin
start() { /etc/webmin/start; }
stop() { /etc/webmin/stop; }